# usr/lib/systemd/user/regolith-x11.target
[Unit]
Description=Regolith on X11
DefaultDependencies=no

Requisite=gnome-session-initialized.target
PartOf=gnome-session-initialized.target
Before=gnome-session-initialized.target
